# rc.local
#! /bin/bash

/etc/scripts/add-my-arp.sh

exit 0
