# rc.local
#! /bin/bash

#/etc/scripts/ipv6.sh
/etc/scripts/add-my-arp.sh
systemctl start firewall.service

exit 0
