# rc.local
#! /bin/bash

/etc/add-my-arp.sh

exit 0
